# packaging/jackfreq.service.in
[Unit]
Description=@PROJECT_DESCRIPTION@

[Service]
Environment=OPTIONS="-w -vv"
EnvironmentFile=-/etc/default/jackfreqd
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/jackfreqd $OPTIONS
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
